# oracleu
#!/bin/bash
# chkconfig: 2345 95 20
# description: Starts the oracle processes on this system
#processname: oracleu
su - oracleu_all
/bin/touch /var/lock/subsys/oracled
