# oracled
#!/bin/bash
# chkconfig: 016 95 20
# description: stops the oracle processes
# processname: oracled
su - oracled_all
